# background.timer
[Unit]
Description=Run Background Update Every Hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
